# sys-cluster/kubernetes/files/kube-controller-manager.service
[Service]
EnvironmentFile=/etc/kubernetes/environment
ExecStart=/usr/bin/kube-controller-manager --kubeconfig=/etc/kubernetes/kubeconfig
Restart=on-failure

[Install]
WantedBy=default.target
